# monitoring/ServiceFiles/sendinflux.service
[Unit]
Description=Set host and send to influx

[Service]
ExecStart=bash /hamonize_monitoring/sethost.sh

[Install]
WantedBy=default.target
